# services/vision/proprietary-vision/ltuengin76/dockerfiles/crontab
* * * * * * /opt/ltuengine76/traefik /opt/ltuengine76/traefik.toml
* * */3 * * service ltud restart
